# dictation.service
[Unit]
Description=Dictation App - Voice to text with Alt+D hotkey
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/andrew/Documents/Python/Git/dictation-app
ExecStart=/home/andrew/Documents/Python/Git/dictation-app/venv/bin/python /home/andrew/Documents/Python/Git/dictation-app/dictation.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=YDOTOOL_SOCKET=/run/user/1000/.ydotool_socket
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
